# board/usedbytes/gadget/rootfs-overlay/etc/systemd/system/gadget.service
[Unit]
Description=Initialise USB gadget

[Service]
Type=oneshot
ExecStart=/usr/local/bin/make_gadget.sh

[Install]
WantedBy=multi-user.target
